# docker/supervisor/conf.d/rabbitmq_consumers.conf
[program:teacher_skills_consumer]
command=php -dmemory_limit=1G -derror_reporting=E_ALL -dlog_errors=1 /app/bin/console rabbitmq:consumer -m 100 teacher_skills --env=consumer -vv
process_name=%(program_name)s_%(process_num)02d
numprocs=2  # Меньше процессов для тяжелых операций
directory=/app
autostart=true
autorestart=true
startsecs=3
startretries=3
stopwaitsecs=10
startdelay=10
user=otus
redirect_stderr=true
stdout_logfile=/app/var/log/supervisor/teacher_skills.log
stdout_capture_maxbytes=1MB
stderr_logfile=/app/var/log/supervisor/teacher_skills.error.log
stderr_capture_maxbytes=1MB

[program:teacher_skills_dlq_consumer]
command=php -dmemory_limit=1G -derror_reporting=E_ALL -dlog_errors=1 /app/bin/console rabbitmq:consumer -m 100 teacher_skills_dlq --env=consumer -vv
process_name=%(program_name)s_%(process_num)02d
numprocs=1  # Для DLQ достаточно одного процесса
directory=/app
autostart=true
autorestart=true
startsecs=3
startretries=3
stopwaitsecs=10
startdelay=10
user=otus
redirect_stderr=true
stdout_logfile=/app/var/log/supervisor/teacher_skills_dlq.log
stdout_capture_maxbytes=1MB
stderr_logfile=/app/var/log/supervisor/teacher_skills_dlq.error.log
stderr_capture_maxbytes=1MB

[program:teacher_groups_consumer]
command=php -dmemory_limit=1G -derror_reporting=E_ALL -dlog_errors=1 /app/bin/console rabbitmq:consumer -m 100 teacher_groups --env=consumer -vv
process_name=%(program_name)s_%(process_num)02d
numprocs=2  # Меньше процессов для тяжелых операций
directory=/app
autostart=true
autorestart=true
startsecs=3
startretries=3
stopwaitsecs=10
startdelay=10
user=otus
redirect_stderr=true
stdout_logfile=/app/var/log/supervisor/teacher_groups.log
stdout_capture_maxbytes=1MB
stderr_logfile=/app/var/log/supervisor/teacher_groups.error.log
stderr_capture_maxbytes=1MB

[program:teacher_groups_dlq_consumer]
command=php -dmemory_limit=1G -derror_reporting=E_ALL -dlog_errors=1 /app/bin/console rabbitmq:consumer -m 100 teacher_groups_dlq --env=consumer -vv
process_name=%(program_name)s_%(process_num)02d
numprocs=1
directory=/app
autostart=true
autorestart=true
startsecs=3
startretries=3
stopwaitsecs=10
startdelay=10
user=otus
redirect_stderr=true
stdout_logfile=/app/var/log/supervisor/teacher_groups_dlq.log
stdout_capture_maxbytes=1MB
stderr_logfile=/app/var/log/supervisor/teacher_groups_dlq.error.log
stderr_capture_maxbytes=1MB

[program:student_skills_consumer]
command=php -dmemory_limit=1G -derror_reporting=E_ALL -dlog_errors=1 /app/bin/console rabbitmq:consumer -m 100 student_skills --env=consumer -vv
process_name=%(program_name)s_%(process_num)02d
numprocs=3  # Средняя нагрузка
directory=/app
autostart=true
autorestart=true
startsecs=3
startretries=3
stopwaitsecs=10
startdelay=10
user=otus
redirect_stderr=true
stdout_logfile=/app/var/log/supervisor/student_skills.log
stdout_capture_maxbytes=1MB
stderr_logfile=/app/var/log/supervisor/student_skills.error.log
stderr_capture_maxbytes=1MB

[program:student_skills_dlq_consumer]
command=php -dmemory_limit=1G -derror_reporting=E_ALL -dlog_errors=1 /app/bin/console rabbitmq:consumer -m 100 student_skills_dlq --env=consumer -vv
process_name=%(program_name)s_%(process_num)02d
numprocs=1
directory=/app
autostart=true
autorestart=true
startsecs=3
startretries=3
stopwaitsecs=10
startdelay=10
user=otus
redirect_stderr=true
stdout_logfile=/app/var/log/supervisor/student_skills_dlq.log
stdout_capture_maxbytes=1MB
stderr_logfile=/app/var/log/supervisor/student_skills_dlq.error.log
stderr_capture_maxbytes=1MB

[program:student_groups_consumer]
command=php -dmemory_limit=1G -derror_reporting=E_ALL -dlog_errors=1 /app/bin/console rabbitmq:consumer -m 100 student_groups --env=consumer -vv
process_name=%(program_name)s_%(process_num)02d
numprocs=3  # Средняя нагрузка
directory=/app
autostart=true
autorestart=true
startsecs=3
startretries=3
stopwaitsecs=10
startdelay=10
user=otus
redirect_stderr=true
stdout_logfile=/app/var/log/supervisor/student_groups.log
stdout_capture_maxbytes=1MB
stderr_logfile=/app/var/log/supervisor/student_groups.error.log
stderr_capture_maxbytes=1MB

[program:student_groups_dlq_consumer]
command=php -dmemory_limit=1G -derror_reporting=E_ALL -dlog_errors=1 /app/bin/console rabbitmq:consumer -m 100 student_groups_dlq --env=consumer -vv
process_name=%(program_name)s_%(process_num)02d
numprocs=1
directory=/app
autostart=true
autorestart=true
startsecs=3
startretries=3
stopwaitsecs=10
startdelay=10
user=otus
redirect_stderr=true
stdout_logfile=/app/var/log/supervisor/student_groups_dlq.log
stdout_capture_maxbytes=1MB
stderr_logfile=/app/var/log/supervisor/student_groups_dlq.error.log
stderr_capture_maxbytes=1MB

[program:group_skills_consumer]
command=php -dmemory_limit=1G -derror_reporting=E_ALL -dlog_errors=1 /app/bin/console rabbitmq:consumer -m 100 group_skills --env=consumer -vv
process_name=%(program_name)s_%(process_num)02d
numprocs=4  # Больше процессов для легких операций
directory=/app
autostart=true
autorestart=true
startsecs=3
startretries=3
stopwaitsecs=10
startdelay=10
user=otus
redirect_stderr=true
stdout_logfile=/app/var/log/supervisor/group_skills.log
stdout_capture_maxbytes=1MB
stderr_logfile=/app/var/log/supervisor/group_skills.error.log
stderr_capture_maxbytes=1MB

[program:group_skills_dlq_consumer]
command=php -dmemory_limit=1G -derror_reporting=E_ALL -dlog_errors=1 /app/bin/console rabbitmq:consumer -m 100 group_skills_dlq --env=consumer -vv
process_name=%(program_name)s_%(process_num)02d
numprocs=1
directory=/app
autostart=true
autorestart=true
startsecs=3
startretries=3
stopwaitsecs=10
startdelay=10
user=otus
redirect_stderr=true
stdout_logfile=/app/var/log/supervisor/group_skills_dlq.log
stdout_capture_maxbytes=1MB
stderr_logfile=/app/var/log/supervisor/group_skills_dlq.error.log
stderr_capture_maxbytes=1MB

[program:cache_invalidation_consumer]
command=php -dmemory_limit=1G -derror_reporting=E_ALL -dlog_errors=1 /app/bin/console rabbitmq:consumer -m 100 cache_invalidation --env=consumer -vv
process_name=%(program_name)s_%(process_num)02d
numprocs=2
directory=/app
autostart=true
autorestart=true
startsecs=3
startretries=3
stopwaitsecs=10
startdelay=10
user=otus
redirect_stderr=true
stdout_logfile=/app/var/log/supervisor/cache_invalidation.log
stdout_capture_maxbytes=1MB
stderr_logfile=/app/var/log/supervisor/cache_invalidation.error.log
stderr_capture_maxbytes=1MB

[program:cache_invalidation_dlq_consumer]
command=php -dmemory_limit=1G -derror_reporting=E_ALL -dlog_errors=1 /app/bin/console rabbitmq:consumer -m 100 cache_invalidation_dlq --env=consumer -vv
process_name=%(program_name)s_%(process_num)02d
numprocs=1
directory=/app
autostart=true
autorestart=true
startsecs=3
startretries=3
stopwaitsecs=10
startdelay=10
user=otus
redirect_stderr=true
stdout_logfile=/app/var/log/supervisor/cache_invalidation_dlq.log
stdout_capture_maxbytes=1MB
stderr_logfile=/app/var/log/supervisor/cache_invalidation_dlq.error.log
stderr_capture_maxbytes=1MB

[program:domain_events_consumer]
command=php -dmemory_limit=1G -derror_reporting=E_ALL -dlog_errors=1 /app/bin/console rabbitmq:consumer -m 100 domain_events --env=consumer -vv
process_name=%(program_name)s_%(process_num)02d
numprocs=2
directory=/app
autostart=true
autorestart=true
startsecs=3
startretries=3
stopwaitsecs=10
startdelay=10
user=otus
redirect_stderr=true
stdout_logfile=/app/var/log/supervisor/domain_events.log
stdout_capture_maxbytes=1MB
stderr_logfile=/app/var/log/supervisor/domain_events.error.log
stderr_capture_maxbytes=1MB

[program:domain_events_dlq_consumer]
command=php -dmemory_limit=1G -derror_reporting=E_ALL -dlog_errors=1 /app/bin/console rabbitmq:consumer -m 100 domain_events_dlq --env=consumer -vv
process_name=%(program_name)s_%(process_num)02d
numprocs=1
directory=/app
autostart=true
autorestart=true
startsecs=3
startretries=3
stopwaitsecs=10
startdelay=10
user=otus
redirect_stderr=true
stdout_logfile=/app/var/log/supervisor/domain_events_dlq.log
stdout_capture_maxbytes=1MB
stderr_logfile=/app/var/log/supervisor/domain_events_dlq.error.log
stderr_capture_maxbytes=1MB
